# devconfig/debadm_MBP14_deCH/home/fortinet/.config/autostart/set-resolution.desktop
[Desktop Entry]
Type=Application
Exec=/home/fortinet/bin/res
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Set Custom Resolution
